# devops/arch/flatscrape.timer
[Unit]
Description=Flat scraping

[Timer]
OnStartupSec=10s
OnUnitInactiveSec=600s

[Install]
WantedBy=multi-user.target
